# server/mc@.socket
[Socket]
ListenFIFO=%t/%N
RemoveOnStop=true
SocketMode=660
SocketUser=mc
SocketGroup=mc

[Install]
WantedBy=sockets.target
